# Lab2/battery-timer.service
[Unit]
Description=Battery timer

[Service]
Type=simple
ExecStart=/home/longdude/battery_timer.sh
User=longdude

[Install]
WantedBy=getty.target
