# Lab2/battery-timer.timer
[Unit]
Description=Battery Analyzer timer

[Timer]
OnUnitActiveSec=30s
Unit=battery-timer.service
AccuracySec=1s

[Install]
WantedBy=timers.target
